# x11-drivers/nvidia-drivers/files/nvidia-vgpud.initd
#!/sbin/openrc-run
# Copyright 1999-2021 Gentoo Authors
# Distributed under the terms of the GNU General Public License v2

description="NVIDIA vGPU Daemon"
command="nvidia-vgpud"
command_args="${NVPD_USER:+--user ${NVPD_USER}} ${ARGS}"
pidfile="/var/run/nvidia-vgpud/nvidia-vgpud.pid"

stop_post() {
	rmdir "${pidfile%/*}" 2>/dev/null || true
}
